# files/etc/uci-defaults/led_setup.sh
#!/bin/sh

sh /etc/init.d/turn_on_power_led.sh enable
sh /etc/init.d/start_wifi_led.sh enable
